# ntopng-extras/backupConfig.service
[Unit]
Description=backing up the ntopng configuration

[Service]
Type=oneshot
ExecStart=/volume1/ntopng-extras/backupConfig

