# ntopng-extras/backupConfig.timer
[Unit]
Description=backing up the ntopng configuration

[Timer]
OnCalendar=Mon *-*-* 01:00:00
Persistent=true

[Install]
WantedBy=timers.target
